# srs-thumbnail/srs-monitor.timer
[Unit]
Description=Timer for IBLUPS SRS Thumbnail Monitor

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
Unit=srs-monitor.service

[Install]
WantedBy=timers.target